# rootfs/etc/s6-overlay/s6-rc.d/redis/run
#!/command/with-contenv bash
install -d -o redis -g redis /data/redis
exec s6-setuidgid redis redis-server --bind 127.0.0.1 --port 6379 --requirepass "${REDIS_PASSWORD}" --dir /data/redis --save ""
