# rootfs/etc/s6-overlay/s6-rc.d/postgres/run
#!/command/with-contenv bash
PG_BIN=$(ls -d /usr/lib/postgresql/*/bin | head -1)
exec s6-setuidgid postgres "$PG_BIN/postgres" -D /data/postgres -c listen_addresses=127.0.0.1 -c port=5432
